# deb/backend/DEBIAN/postrm
#!/bin/sh -e

update-rc.d -f openbook-backend remove
