# _deb/frontend/etc/cron.d/openbook-subscription-notifyer
0 2 * * * openbook /usr/lib/openbook/utils/subscription-notifyer -config /etc/openbook/frontend/config.gcfg
