# _build/frontend/DEBIAN/postrm
#!/bin/sh -e

update-rc.d -f openbook-frontend remove